# deploy/cashflux-retention.example.service
[Unit]
Description=CashFlux self-host retention pruning
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/cashflux
ExecStart=/usr/bin/docker compose -f docker-compose.selfhost.yml run --rm cashflux-server retention
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
